# feishu-script-web/feishu-script.service
[Unit]
Description=Feishu Script Platform
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/app
ExecStart=/usr/bin/node .output/server/index.mjs
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3030
Environment=NUXT_FEISHU_APP_ID=your_app_id
Environment=NUXT_FEISHU_APP_SECRET=your_app_secret

# 日志配置
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=feishu-script

[Install]
WantedBy=multi-user.target
